# serf-agent.conf
# Serf 0.2.x Agent (Upstart unit)
description "Serf Agent 1"
start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [06]

kill signal INT    # Use SIGINT instead of SIGTERM so serf can depart the cluster.
respawn            # Restart the process if it dies and GOAL was not 'stopping'.
kill timeout 90    # Allow 90 seconds for serf to die before sending SIGKILL.

env SERF=/usr/local/bin/serf

exec $SERF agent -iface=eth1 -rpc-addr=0.0.0.0:7373
